# images/raiden/0.100.5a1.dev162-2217bcb-simnet/supervisord.conf
[supervisord]
nodaemon=true
logfile=/root/supervisord.log
childlogdir=/root
user=root

[program:raiden]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
command=/opt/entrypoint.sh
